# scripts/openterminalui.service
[Unit]
Description=OpenTerminalUI — веб-терминал для трейдеров
After=network.target

[Service]
Type=simple
User={{APP_USER}}
Group={{APP_GROUP}}
WorkingDirectory={{INSTALL_DIR}}
Environment="PYTHONPATH={{INSTALL_DIR}}"
EnvironmentFile=-{{INSTALL_DIR}}/.env
ExecStart={{INSTALL_DIR}}/.venv/bin/python -m uvicorn backend.main:app --host 0.0.0.0 --port {{PORT}}
Restart=on-failure
RestartSec=5
StandardOutput=append:{{LOG_DIR}}/openterminalui.log
StandardError=append:{{LOG_DIR}}/openterminalui.log

[Install]
WantedBy=multi-user.target
